# linux systemctl Services/moneyServerDBbackuptask.service
[Unit]
Description=moneyServerDBbackupTask

[Service]
ExecStart=bash /home/parth/serverUtilsForMoneyServer/moneyServerBKP.sh
